# deployment/bella-backend.service
[Unit]
Description=Bella Subscription Platform Backend API
After=network.target postgresql.service redis.service
Requires=postgresql.service redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/bella/backend
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/var/www/bella/backend"
ExecStart=/var/www/bella/backend/venv/bin/python -m uvicorn src.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bella-backend

[Install]
WantedBy=multi-user.target
