# deployment/bella-scheduler.service
[Unit]
Description=Bella Subscription Platform Background Scheduler
After=network.target bella-backend.service
Requires=bella-backend.service

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/bella/backend
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/var/www/bella/backend"
ExecStart=/var/www/bella/backend/venv/bin/python /var/www/bella/backend/scripts/run_scheduler.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bella-scheduler

[Install]
WantedBy=multi-user.target
